# deploy/trpg-master-staging-backup.timer
[Unit]
Description=Daily encrypted TRPG Master staging backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=30m
Unit=trpg-master-staging-backup.service

[Install]
WantedBy=timers.target
