# deploy/trpg-master-staging-backup.service
[Unit]
Description=Encrypted TRPG Master staging database and runtime backup
After=postgresql.service

[Service]
Type=oneshot
User=trpgdeploy
Group=trpgdeploy
EnvironmentFile=/etc/trpg-master/staging.env
Environment=TRPG_BACKUP_ROOT=/var/backups/trpg-master-staging
Environment=TRPG_BACKUP_RUNTIME_ROOT=/var/lib/trpg-master-staging
Environment=TRPG_BACKUP_PREFIX=trpg-master-staging
ExecStart=/opt/trpg-master-staging/current/deploy/backup-trpg-master.sh
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/trpg-master-staging
